# supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /efit-institute-api/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx